# overlay/etc/systemd/system/docker.service.d/override.conf
[Service]
ExecStart=
ExecStart=/usr/bin/docker daemon --storage-driver overlay2 -H fd:// $DOCKER_OPTS
