# overlay/etc/systemd/system/ufw.service.d/override.conf
[Service]
ExecStartPre=/usr/sbin/ufw logging off
